# experiments/config/publisher.service
[Unit]
Description=Spectrum publisher

[Service]
ExecStart=/bin/bash -c "/home/ubuntu/spectrum/publisher --local-port 6001 --public-address $(hostname):5001 --log-level trace"
Type=simple
Restart=no
EnvironmentFile=/etc/spectrum.conf
Environment="RUST_BACKTRACE=1"
LimitNOFILE=64000
LimitFSIZE=infinity
LimitCPU=infinity
LimitAS=infinity
LimitMEMLOCK=infinity
LimitNPROC=64000
