# src/experiment/data/spectrum.service.template
[Unit]
Description=Spectrum ${service}

[Service]
ExecStart=${binary} --local-port ${port} --public-address ${public_address} --log-level info
Type=simple
Restart=no
EnvironmentFile=/etc/spectrum.conf
LimitNOFILE=64000
LimitFSIZE=infinity
LimitCPU=infinity
LimitAS=infinity
LimitMEMLOCK=infinity
LimitNPROC=64000
