# experiments_new/config/leader.service
[Unit]
Description=Spectrum leader

[Service]
ExecStart=/bin/bash -c "/home/ubuntu/spectrum/leader --local-port 6000 --public-port $(hostname):5000 --log-level trace"
Type=simple
Restart=no
EnvironmentFile=/etc/spectrum.conf
Environment="RUST_BACKTRACE=1"
LimitNOFILE=64000
LimitFSIZE=infinity
LimitCPU=infinity
LimitAS=infinity
LimitMEMLOCK=infinity
LimitNPROC=64000
