# experiments/spectrum/config/broadcaster@.service
# Start as:
#    systemctl broadcaster{0..10}.service
# Then, everywhere you see "%i" in this file will be replaced by the client index.

[Unit]
Description=Spectrum client

[Service]
ExecStart=/home/ubuntu/spectrum/broadcaster --index %i
Type=simple
Restart=no
EnvironmentFile=/etc/spectrum.conf
