# experiments/spectrum/config/viewer@.service
# Start as:
#    systemctl viewer{300..399}.service

[Unit]
Description=Spectrum client

[Service]
ExecStart=/home/ubuntu/spectrum/viewer
Type=simple
Restart=no
EnvironmentFile=/etc/spectrum.conf
Environment="RUST_BACKTRACE=1"
