# experiments/spectrum/config/worker@.service
[Unit]
Description=Spectrum worker %i

[Service]
# Bash is needed for arithmetic with %i
ExecStart=/bin/bash -c "/home/ubuntu/spectrum/worker --local-port $((6100 + %i - 1)) --public-address $(hostname):$((5100 + %i - 1)) --index $(($SPECTRUM_WORKER_START_INDEX + %i)) --log-level trace"
Type=simple
Restart=no
EnvironmentFile=/etc/spectrum.conf
Environment="RUST_BACKTRACE=1"
LimitNOFILE=64000
LimitFSIZE=infinity
LimitCPU=infinity
LimitAS=infinity
LimitMEMLOCK=infinity
LimitNPROC=64000
